# contrib/packaging/debian/postinst
#!/bin/sh
set -e

echo ""
echo "To enable the daemon as a user service:"
echo "  systemctl --user daemon-reload"
echo "  systemctl --user enable --now agent-hub-daemon"
echo ""
echo "Or use the built-in installer:"
echo "  agent-hub-daemon --install-service"
echo ""

#DEBHELPER#

exit 0
